# roles/docker_prune/files/docker_prune.service
[Unit]
Description=Docker Prune

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune --all --volumes --force
WorkingDirectory=/tmp
