# roles/autorestic/templates/autorestic_backup.timer.j2
[Unit]
Description=Backup {{ autorestic_location }} timer

[Timer]
OnCalendar=*-*-* 12:00:00
Unit=autorestic_backup_{{ autorestic_location }}.service

[Install]
WantedBy=timers.target
