# roles/autorestic/templates/autorestic_backup.service.j2
[Unit]
Description=Backup {{ autorestic_location }} using autorestic

[Service]
Type=simple
User={{ ansible_ssh_user | default(ansible_user_id) }}
WorkingDirectory={{ service_path }}
ExecStart=docker compose run autorestic bash -c "apk add curl && autorestic backup --verbose --location {{ autorestic_location }}"
